# retrophone.service
[Unit]
Description=RetroPhone Time Travel Radio
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/RetroPhone
ExecStart=/bin/bash /home/pi/RetroPhone/start.sh
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
